# .config/systemd/user/mbsync-prio-0.timer
[Unit]
Description=mbsync query timer [prio-0]
ConditionPathExists=%h/.mbsyncrc

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1m

[Install]
WantedBy=default.target
